# deployment/vps/nyiso-api.service
[Unit]
Description=NYISO FastAPI Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/nyiso-dashboard
Environment="PATH=/opt/nyiso-dashboard/venv/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="API_HOST=127.0.0.1"
Environment="API_PORT=8000"
Environment="DATABASE_URL=sqlite:///opt/nyiso-dashboard/nyiso_data.db"
ExecStart=/opt/nyiso-dashboard/venv/bin/python start_api_prod.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nyiso-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nyiso-dashboard

[Install]
WantedBy=multi-user.target


