# systemd/nyiso-api.service
[Unit]
Description=NYISO FastAPI Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/nyiso-dashboard
Environment="PATH=/opt/nyiso-dashboard/venv/bin"
Environment="DATABASE_URL=sqlite:///opt/nyiso-dashboard/nyiso_data.db"
ExecStart=/opt/nyiso-dashboard/venv/bin/python /opt/nyiso-dashboard/start_api_prod.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

